# deploy/listen.conf
#!upstart

start on runlevel [23]

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

script
  chdir /home/ubuntu
  exec ./listen 2>&1
end script
